# chef/cookbooks/ironic/templates/default/tftp.service.erb
[Unit]
Description=Tftp Server

[Service]
Type=simple
ExecStart=/usr/sbin/in.tftpd -u tftp -s <%= @tftproot %> -m <%= @map_file %> -L -a <%= @ironic_ip %>
