# chef/cookbooks/cinder/templates/default/cinder-looplvm.service.erb
[Unit]
Description=Loopback mount for LVM VGs used by cinder
Before=<%= @service_name %>.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/cinder-looplvm start
ExecReload=/usr/bin/cinder-looplvm start
ExecStop=/usr/bin/cinder-looplvm stop

[Install]
WantedBy=multi-user.target
