# chef/cookbooks/monasca/templates/default/elasticsearch-curator-systemd.service.erb
[Unit]
Description=Run elasticsearch-curator

[Service]
Type=oneshot
User=curator
Group=curator
ExecStart=/usr/bin/curator --config /etc/elasticsearch-curator/curator_config.yml /etc/elasticsearch-curator/curator_action.yml
