# chef/cookbooks/monasca/templates/default/elasticsearch-curator-systemd.timer.erb
[Unit]
Description=Run elasticsearch-curator

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target